# crontab
*/5 * * * * root sh /home/osboxes/add_record_script.sh
